# srv/salt/haproxy/haproxy-acme-dump.timer
[Unit]
Description=Persist HAProxy ACME certificates regularly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
